# jasanjejop/jasanjejop-bot.service
[Unit]
Description=자산제곱 분석 봇
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-bots/jasanjejop
Environment=PATH=/home/ubuntu/telegram-bots/jasanjejop/venv/bin
ExecStart=/home/ubuntu/telegram-bots/jasanjejop/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
